# deploy/shoptest-bug-maintenance.timer
[Unit]
Description=Run ShopTest submitted BUG maintenance every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true
Unit=shoptest-bug-maintenance.service

[Install]
WantedBy=timers.target
